# build/desktops/app-profiles/appconfs/blender.conf
[program:blender]
command=/usr/local/blender-3.6.5-linux-%ARCH%/blender
environment=DISPLAY=":10",PULSE_SERVER="/var/run/kvdi/pulse-server",HOME="/var/run/kvdi"
autostart=true
autorestart=true
startsecs=0
startretries=1000000
priority=99
