# build/desktops/dosbox/supervisor/display.conf
[program:display]
command=/usr/bin/Xvnc :10 -rfbunixpath /var/run/kvdi/display.sock -SecurityTypes None
autostart=true
autorestart=true
startsecs=5
startretries=1000000