# build/desktops/dosbox/supervisor/dosbox.conf
[program:dosbox]
command=/usr/bin/dosbox -conf /dos/dosbox.conf
environment=DISPLAY=":10"
autostart=true
autorestart=true
startsecs=5
startretries=1000000