# build/desktops/dosbox/supervisor/dosbox.conf
[program:dosbox]
command=/bin/sh -c 'sleep 3 && /usr/bin/dosbox -conf /dos/dosbox.conf'
environment=DISPLAY=":10",SDL_AUDIODRIVER="pulse",PULSE_SERVER="/var/run/kvdi/pulse-server",HOME="/var/run/kvdi"
autostart=true
autorestart=true
startsecs=0
startretries=1000000
priority=99