# build/desktops/dosbox/supervisor/pulseaudio.conf
[program:pulseaudio]
command=/usr/bin/pulseaudio --realtime=true 
environment=DISPLAY=":10"
autostart=true
autorestart=true
startsecs=5
startretries=1000000