# build/desktops/app-profiles/rootfs/etc/systemd/system/user-init.service
[Unit]
Description=SystemD User Init Process
Wants=network.target

[Service]
Type=idle
ExecStart=/usr/local/sbin/fakegetty
EnvironmentFile=/etc/default/display
UtmpIdentifier=cons
TTYVHangup=yes
KillMode=process
IgnoreSIGPIPE=no
SendSIGHUP=yes
Restart=always

[Install]
WantedBy=multi-user.target