# build/desktops/app-profiles/systemd/Lens.service
[Unit]
Description=Lens
After=display.service
Requires=display.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/Lens
EnvironmentFile=/etc/default/desktop

[Install]
WantedBy=default.target
