# build/desktops/app-profiles/systemd/Postman.service
[Unit]
Description=Postman
After=display.service
Requires=display.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/Postman
EnvironmentFile=/etc/default/desktop

[Install]
WantedBy=default.target

