# build/desktops/arch/rootfs/etc/systemd/user/display.service
[Unit]
Description=Xvnc Display

[Service]
Type=simple
Restart=always
ExecStart=/usr/sbin/Xvnc ${DISPLAY} -rfbunixpath ${VNC_SOCK_ADDR} -SecurityTypes None
EnvironmentFile=/etc/default/xvnc

[Install]
WantedBy=default.target
