# build/desktops/app-profiles/rootfs/etc/systemd/user/display.service
[Unit]
Description=Xpra Display

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/xpra --no-daemon --bind=${VNC_SOCK_ADDR} start ${DISPLAY}
EnvironmentFile=/etc/default/kvdi

[Install]
WantedBy=default.target
