# build/desktops/arch/systemd/plasma.service
[Unit]
Description=Plasma
After=display.service
Requires=display.service

[Service]
Type=simple
Restart=always
ExecStart=/bin/bash -c 'export $$(dbus-launch) ; exec /usr/sbin/startplasma-x11'
EnvironmentFile=/etc/default/kvdi

[Install]
WantedBy=default.target
