# build/desktops/arch/systemd/xfce.service
[Unit]
Description=Plasma
After=display.service
Requires=display.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/sbin/startxfce4
EnvironmentFile=/etc/default/kvdi

[Install]
WantedBy=default.target
