# build/desktops/app-profiles/lens.conf
[program:lens]
command=/bin/sh -c 'cd /var/run/kvdi && /usr/local/bin/Lens --appimage-extract && /var/run/kvdi/squashfs-root/AppRun'
environment=DISPLAY=":10",PULSE_SERVER="/var/run/kvdi/pulse-server",HOME="/var/run/kvdi"
autostart=true
autorestart=true
startsecs=0
startretries=1000000
priority=99