# build/desktops/app-profiles/supervisor/lwm.conf
[program:lwm]
command=/usr/bin/lwm
environment=DISPLAY=":10",HOME="/var/run/kvdi"
autostart=true
autorestart=true
startsecs=5
startretries=1000000
priority=2