# build/desktops/app-profiles/supervisor/pulseaudio.conf
[program:pulseaudio]
command=/usr/bin/pulseaudio
environment=DISPLAY=":10",HOME="/var/run/kvdi"
autostart=true
autorestart=true
startsecs=3
startretries=1000000
priority=2