# build/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
redirect_stderr=true
stopsignal=TERM
autorestart=true
directory=/root
user=root
childlogdir=/var/log/supervisord

[group:x]
programs=xvnc

[program:xvnc]
priority=1
command=Xvnc :10 -rfbunixpath %UNIX_SOCK% -SecurityTypes None
